# build/DEBIAN/postinst
#!/bin/sh
# postinst script for transmission
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

get_share() {
    TRANSMISSION_SHARE=`rn_nml -g shares | grep Transmission | cut -d\" -f2`;
}


case "$1" in
    configure)
    sysctl -w net.core.rmem_max=4194304
    sysctl -w net.core.wmem_max=1048576
    
    get_share

    if [ -z "$TRANSMISSION_SHARE" ]
    then

        # create Transmission share 
        rn_nml -a 'share:/*/Transmission'

        get_share

        mkdir /$TRANSMISSION_SHARE/watch-dir
        mkdir /$TRANSMISSION_SHARE/incomplete

    fi


    if [ ! -r "/$TRANSMISSION_SHARE/settings.json" ]
    then
        cp /apps/transmission/settings.json.orig /$TRANSMISSION_SHARE/settings.json

        echo "    "\"watch-dir\": \"/$TRANSMISSION_SHARE/watch-dir/\", >> /$TRANSMISSION_SHARE/settings.json
        echo "    "\"download-dir\": \"/$TRANSMISSION_SHARE/downloads\", >> /$TRANSMISSION_SHARE/settings.json
        echo "    "\"incomplete-dir\": \"/$TRANSMISSION_SHARE/incomplete\" >> /$TRANSMISSION_SHARE/settings.json
        echo "}" >> /$TRANSMISSION_SHARE/settings.json
    fi

    echo TRANSMISSION_HOME=/$TRANSMISSION_SHARE > /apps/transmission/fvapp-transmission.env

    chown -R admin:admin /apps/transmission
    chown -R guest:guest /$TRANSMISSION_SHARE/

    if systemctl restart apache2.service; then
      # success
      event_push app readynasd '<add-s resource-type="LocalApp" resource-id="LocalApp"><LocalApp appname="transmission" success="1"/></add-s>' 0 0
    else
      # error
      event_push app readynasd '<add-s resource-type="LocalApp" resource-id="LocalApp"><LocalApp appname="transmission" success="0"/></add-s>' 0 0
    fi

    systemctl enable fvapp-transmission.service
    systemctl start fvapp-transmission.service
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.



exit 0
